# contrib/init/wescoind.conf
description "Wescoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env WESCOIND_BIN="/usr/bin/wescoind"
env WESCOIND_USER="wescoin"
env WESCOIND_GROUP="wescoin"
env WESCOIND_PIDDIR="/var/run/wescoind"
# upstart can't handle variables constructed with other variables
env WESCOIND_PIDFILE="/var/run/wescoind/wescoind.pid"
env WESCOIND_CONFIGFILE="/etc/wescoin/wescoin.conf"
env WESCOIND_DATADIR="/var/lib/wescoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # wescoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$WESCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run wescoind."
        echo "The setting must appear in $WESCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Wescoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$WESCOIND_PIDDIR"
    chmod 0755 "$WESCOIND_PIDDIR"
    chown $WESCOIND_USER:$WESCOIND_GROUP "$WESCOIND_PIDDIR"
    chown $WESCOIND_USER:$WESCOIND_GROUP "$WESCOIND_CONFIGFILE"
    chmod 0660 "$WESCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$WESCOIND_PIDFILE" \
    --chuid $WESCOIND_USER:$WESCOIND_GROUP \
    --exec "$WESCOIND_BIN" \
    -- \
    -pid="$WESCOIND_PIDFILE" \
    -conf="$WESCOIND_CONFIGFILE" \
    -datadir="$WESCOIND_DATADIR" \
    -disablewallet \
    -daemon

